# bb_biobase_local/mdns/biobase-mdns.service
# Publishes biobase.local (A record) for the LAN. Does not modify avahi-daemon.conf.
# Install:  sudo ./install-biobase-mdns.sh
# Manual:   sudo install -m0755 publish-biobase-local.sh /usr/local/bin/biobase-mdns-publish
#           sudo install -m0644 biobase-mdns.service /etc/systemd/system/
#           sudo systemctl daemon-reload && sudo systemctl enable --now biobase-mdns.service

[Unit]
Description=Publish biobase.local (mDNS address record for Biobase)
After=network-online.target avahi-daemon.service
Wants=network-online.target
Requires=avahi-daemon.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/biobase-mdns
ExecStart=/usr/local/bin/biobase-mdns-publish
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
